# asa_group_sync.timer
[Unit]
Description=Run asa_group_sync.service every 5 minutes

[Timer]
OnCalendar=*:0/5
Persistent=true

[Install]
WantedBy=timers.target
